# scripts/datasync/datasync.timer
[Unit]
Description=Daily synchronization of data partition to its mirror

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
